# systemd/stephenasatsa-frontend.service
[Unit]
Description=Stephen Asatsa Website Frontend (Next.js)
After=network.target stephenasatsa-backend.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/codecrafter/Documents/combined/website
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=append:/home/codecrafter/Documents/combined/website/logs/frontend.log
StandardError=append:/home/codecrafter/Documents/combined/website/logs/frontend-error.log

[Install]
WantedBy=multi-user.target
